# nlddcd.service.in
[Unit]
Description=Netlink-based Dynamic DNS Client Daemon
Documentation=man:nlddcd(8)
After=network.target

[Service]
Type=simple
EnvironmentFile=-@sysconfdir@/sysconfig/nlddcd
ExecStart=@sbindir@/nlddcd $OPTIONS
User=nlddcd
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
